# youtube-api.service
[Unit]
Description=YouTube API Handler Service
After=network.target
Wants=network.target

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/opt/youtube-api
Environment=PATH=/opt/youtube-api/venv/bin
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=False
EnvironmentFile=/opt/youtube-api/production.env
ExecStart=/opt/youtube-api/venv/bin/gunicorn --config gunicorn.conf.py --daemon api_server:create_app()
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
PIDFile=/tmp/youtube_api.pid
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/youtube-api/logs /tmp

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target 